# deploy/systemd/sales-dashboard-backend.service
[Unit]
Description=Sales Dashboard Backend API
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/sales_dashbord/backend
EnvironmentFile=/srv/sales_dashbord/backend/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
